# deploy/gridaibot.service
[Unit]
Description=GridAI Discord Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/gridaibot
Environment=PATH=/opt/gridaibot/.venv/bin:/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/opt/gridaibot/.env
ExecStart=/opt/gridaibot/.venv/bin/python discord_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
